# luci-app-lpac/root/etc/uci-defaults/90-luci-lpac
#!/bin/sh

# This script ONLY runs on FIRST INSTALL
# On upgrades, it does NOTHING to preserve user settings

# Check if this is a first-time install
MARKER_FILE="/etc/luci-lpac/.configured"

# If marker exists, this is an upgrade - EXIT IMMEDIATELY
if [ -f "$MARKER_FILE" ]; then
	# Upgrade detected - do nothing to preserve user settings
	exit 0
fi

# FIRST INSTALL ONLY - from this point forward
mkdir -p /etc/luci-lpac

# Initialize UCI configuration sections
uci -q get luci-lpac.config >/dev/null || uci set luci-lpac.config=luci_lpac
uci -q get luci-lpac.advanced >/dev/null || uci set luci-lpac.advanced=luci_lpac

# Set default values (FIRST INSTALL ONLY)
uci set luci-lpac.config.apdu_driver='qmi'
uci set luci-lpac.config.http_driver='curl'
uci set luci-lpac.config.custom_aid=''
uci set luci-lpac.config.es10x_mss='60'
uci set luci-lpac.config.debug_http='0'
uci set luci-lpac.config.debug_apdu='0'
uci set luci-lpac.config.auto_notification='1'
uci set luci-lpac.config.qmi_slot='0'
uci set luci-lpac.config.qmi_device=''
uci set luci-lpac.config.pcsc_reader=''

uci set luci-lpac.advanced.log_level='info'
uci set luci-lpac.advanced.timeout='120'
uci set luci-lpac.advanced.download_cooldown='60'
uci set luci-lpac.advanced.auto_manage_wwan='1'
uci set luci-lpac.advanced.wwan_interface='wwan'
uci set luci-lpac.advanced.auto_sim_power_cycle='1'

# Commit changes
uci commit luci-lpac

# Mark as configured (prevents script from running on upgrades)
touch "$MARKER_FILE"

# Clear LuCI cache
rm -rf /tmp/luci-*

exit 0
