# luci-app-lpac/root/etc/uci-defaults/90-luci-lpac
#!/bin/sh

# Initialize UCI configuration if not exists
uci -q get luci-lpac.config || {
	uci set luci-lpac.config=luci_lpac
	uci set luci-lpac.config.apdu_driver='pcsc'
	uci set luci-lpac.config.http_driver='curl'
	uci set luci-lpac.config.custom_aid=''
	uci set luci-lpac.config.es10x_mss='60'
	uci set luci-lpac.config.debug_http='0'
	uci set luci-lpac.config.debug_apdu='0'
	uci set luci-lpac.config.auto_notification='1'
	uci set luci-lpac.config.qmi_slot='1'
	uci set luci-lpac.config.pcsc_reader=''
	uci commit luci-lpac
}

# Initialize advanced configuration if not exists
uci -q get luci-lpac.advanced || {
	uci set luci-lpac.advanced=luci_lpac
	uci set luci-lpac.advanced.log_level='info'
	uci set luci-lpac.advanced.timeout='120'
	uci set luci-lpac.advanced.download_cooldown='60'
	uci commit luci-lpac
}

# Clear LuCI cache
rm -rf /tmp/luci-*

exit 0
